# digest/config/launchd/com.hh-research.daily-pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hh-research.daily-pipeline</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/haolinguo/claude code/HH research/daily-digest/scripts/run_daily_pipeline.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/haolinguo/claude code/HH research/daily-digest</string>

    <key>StandardOutPath</key>
    <string>/Users/haolinguo/claude code/HH research/daily-digest/data/logs/cron/launchd_daily.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/haolinguo/claude code/HH research/daily-digest/data/logs/cron/launchd_daily.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>TZ</key>
        <string>Asia/Shanghai</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>

    <!-- If laptop was asleep at fire time, run on wake. -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
